# systemd/.config/systemd/user/low-battery-warning.service
[Unit]
Description=Low Battery Warning

[Service]
ExecStart=$XDG_DATA_HOME/systemd-scripts/low_battery_warning.sh

